# deploy/tsa-kpi-refresh.service
[Unit]
Description=TSA KPI Dashboard data refresh + rebuild (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/tsa-kpi
EnvironmentFile=-/opt/tsa-kpi/.env
Environment=KPI_OUTPUT_DIR=/opt/tsa-kpi/output
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/tsa-kpi/.venv/bin/python /opt/tsa-kpi/scripts/kpi/orchestrate.py
TimeoutStartSec=15min
